# root/etc/services.d/sshd/run
#!/usr/bin/with-contenv bash

exec 2>&1 \
  /usr/sbin/sshd \
    -D \
    -e \
    -p 2222
